# splunk_rules/rule_13_registry_persistence.json
{
  "rule_id": "SPL-013",
  "rule_name": "Persistence - Windows Registry Run Key Modification",
  "severity": "Medium",
  "search_query": "index=windows source=\"XmlWinEventLog:Microsoft-Windows-Sysmon/Operational\" EventCode=13 (TargetObject=\"*\\\\CurrentVersion\\\\Run*\" OR TargetObject=\"*\\\\CurrentVersion\\\\RunOnce*\") | table _time, Computer, User, Image, TargetObject, Details | rename Image as \"Source_Process\", Details as \"Payload_Path\"",
  "description": "Detects modifications to Windows Registry Run or RunOnce keys. Attackers use these keys to ensure their malicious payloads execute automatically at system startup or user logon.",
  "mitre_id": "T1547.001",
  "category": "Persistence",
  "reference_source": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/registry/registry_set_asep_reg_keys_modification.yml",
  "status": "Production"
}
